# qbittorrent-wireguard/root/etc/s6-overlay/s6-rc.d/init-adduser/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

export PUID=$(bashio::config 'PUID')
export PGID=$(bashio::config 'PGID')

PUID=${PUID:-911}
PGID=${PGID:-911}


USERHOME=$(grep abc /etc/passwd | cut -d ":" -f6)
usermod -d "/root" abc

groupmod -o -g "${PGID}" abc
usermod -o -u "${PUID}" abc
usermod -d "${USERHOME}" abc

bashio::log.info "───────────────────────────────────────
User UID:    $(id -u abc)
User GID:    $(id -g abc)
───────────────────────────────────────"

lsiown abc:abc /app
lsiown abc:abc /config
lsiown abc:abc /defaults
