# qbittorrent-wireguard/root/etc/s6-overlay/s6-rc.d/svc-monitor/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

openvpn_enabled=false
wireguard_enabled=false

if bashio::config.true 'openvpn_enabled' ; then
  openvpn_enabled=true
fi
if bashio::config.true 'wireguard_enabled'; then
  wireguard_enabled=true
fi

bashio::log.info "*** starting monitor service ***"

export PATH="/usr/local/sbin:/usr/local/bin:${PATH}"

REAL_IP_FILEv4="/currentipv4"
REAL_IP_FILEv6="/currentipv6"
VPN_CHECK_INTERVAL="${VPN_CHECK_INTERVAL:-300}"
VPN_INFO_URL="${VPN_INFO_URL:-https://ipinfo.io}"

# -------------------------------
# Helpers
# -------------------------------

_fetch_public_ipv4() {
    local resp
    local url
    local curl_iface_opts=()
    local urls=(
        "https://icanhazip.com"
        "https://ifconfig.me/ip"
        "https://api64.ipify.org"
        "https://ipinfo.io/ip"
        "ipecho.net/plain"
    )
    local shuffled_urls

    if [[ -n "${VPN_INTERFACE:-}" ]]; then
        curl_iface_opts+=(--interface "${VPN_INTERFACE}")
    fi

    mapfile -t shuffled_urls < <(printf "%s\n" "${urls[@]}" | shuf)

    for url in "${shuffled_urls[@]}"; do
        resp=$(curl -4 -fsS --max-time 5 "${curl_iface_opts[@]}" "${url}" 2>/dev/null || true)
        resp="${resp//[[:space:]]/}"

        if [[ "${resp}" =~ ^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$ ]] || [[ "${resp}" =~ ^[0-9a-fA-F:]+$ ]]; then
            printf '%s\n' "${resp}"
            return 0
        fi
    done

    return 1
}

_fetch_public_ipv6() {
    local resp
    local url
    local curl_iface_opts=()
    local urls=(
        "https://icanhazip.com"
        "https://ifconfig.me/ip"
        "https://api64.ipify.org"
        "https://ipinfo.io/ip"
        "ipecho.net/plain"
    )
    local shuffled_urls

    if [[ -n "${VPN_INTERFACE:-}" ]]; then
        curl_iface_opts+=(--interface "${VPN_INTERFACE}")
    fi

    mapfile -t shuffled_urls < <(printf "%s\n" "${urls[@]}" | shuf)

    for url in "${shuffled_urls[@]}"; do
        resp=$(curl -6 -fsS --max-time 5 "${curl_iface_opts[@]}" "${url}" 2>/dev/null || true)
        resp="${resp//[[:space:]]/}"

        if [[ "${resp}" =~ ^[0-9]+\.[0-9]+\.[0-9]+\.[0-9]+$ ]] || [[ "${resp}" =~ ^[0-9a-fA-F:]+$ ]]; then
            printf '%s\n' "${resp}"
            return 0
        fi
    done

    return 1
}

read_real_ipv4() {
    local attempt ip

    for attempt in {1..6}; do
        if [[ -f "${REAL_IP_FILEv4}" ]]; then
            ip="$(tr -d '[:space:]' < "${REAL_IP_FILEv4}")"
            if [[ -n "${ip}" ]]; then
                echo "${ip}"
                return 0
            fi
        fi
        sleep 5
    done

    echo ""
    return 0
}

read_real_ipv6() {
    local attempt ip

    for attempt in {1..6}; do
        if [[ -f "${REAL_IP_FILEv6}" ]]; then
            ip="$(tr -d '[:space:]' < "${REAL_IP_FILEv6}")"
            if [[ -n "${ip}" ]]; then
                echo "${ip}"
                return 0
            fi
        fi
        sleep 5
    done

    echo ""
    return 0
}

get_ip_infov4() {
    local ip country json info_url url

    if ! ip="$(_fetch_public_ipv4)"; then
        bashio::log.warning "Unable to determine external IP from fallback IP services."
        return 1
    fi

    country="Unknown"
    info_url="${VPN_INFO_URL:-}"

    if [[ -n "${info_url}" ]]; then
        if [[ "${info_url}" == *"%s"* ]]; then
            printf -v url "${info_url}" "${ip}"
        else
            url="${info_url%/}/${ip}/json"
        fi

        if json="$(curl -4 -fsS --max-time 10 "${url}" 2>/dev/null || true)" && [[ -n "${json}" ]]; then
            country="$(
                printf '%s\n' "${json}" \
                    | sed -n 's/.*"country"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p' \
                    | head -n1
            )"
            [[ -z "${country}" ]] && country="Unknown"
        fi
    fi

    printf '%s %s\n' "${ip}" "${country}"
}

get_ip_infov6() {
    local ip country json info_url url

    if ! ip="$(_fetch_public_ipv6)"; then
        bashio::log.warning "Unable to determine external IP from fallback IP services."
        return 1
    fi

    country="Unknown"
    info_url="${VPN_INFO_URL:-}"

    if [[ -n "${info_url}" ]]; then
        if [[ "${info_url}" == *"%s"* ]]; then
            printf -v url "${info_url}" "${ip}"
        else
            url="${info_url%/}/${ip}/json"
        fi

        if json="$(curl -6 -fsS --max-time 10 "${url}" 2>/dev/null || true)" && [[ -n "${json}" ]]; then
            country="$(
                printf '%s\n' "${json}" \
                    | sed -n 's/.*"country"[[:space:]]*:[[:space:]]*"\([^"]*\)".*/\1/p' \
                    | head -n1
            )"
            [[ -z "${country}" ]] && country="Unknown"
        fi
    fi

    printf '%s %s\n' "${ip}" "${country}"
}

wait_for_vpn_ipv4() {
    local attempt out ip country

    for attempt in {1..5}; do
        if out="$(get_ip_infov4)"; then
            read -r ip country <<< "${out}"

            if [[ -n "${REAL_IPv4:-}" ]] && [[ "${ip}" == "${REAL_IPv4}" ]]; then
                bashio::log.warning \
                    "External IP still equals real IP (${ip}); VPN not ready yet (attempt ${attempt}/5)."
            else
                printf '%s %s\n' "${ip}" "${country}"
                return 0
            fi
        else
            bashio::log.warning "Unable to query external IP (attempt ${attempt}/5)."
        fi

        sleep 5
    done

    return 1
}

wait_for_vpn_ipv6() {
    local attempt out ip country

    for attempt in {1..5}; do
        if out="$(get_ip_infov6)"; then
            read -r ip country <<< "${out}"

            if [[ -n "${REAL_IPv6:-}" ]] && [[ "${ip}" == "${REAL_IPv6}" ]]; then
                bashio::log.warning \
                    "External IP still equals real IP (${ip}); VPN not ready yet (attempt ${attempt}/5)."
            else
                printf '%s %s\n' "${ip}" "${country}"
                return 0
            fi
        else
            bashio::log.warning "Unable to query external IP (attempt ${attempt}/5)."
        fi

        sleep 5
    done

    return 1
}

# -------------------------------
# Main logic
# -------------------------------

vpn_openvpn=false
vpn_wireguard=false

if bashio::config.true 'openvpn_enabled'; then
    vpn_openvpn=true
fi

if [[ "${vpn_openvpn}" == true ]] && ! bashio::config.true 'openvpn_alt_mode'; then
    VPN_INTERFACE="tun0"
    bashio::log.info "VPN monitor set to query external IP through interface ${VPN_INTERFACE} (interface binding)."
else
    VPN_INTERFACE=""
fi

if bashio::config.true 'wireguard_enabled'; then
    vpn_wireguard=true
fi

if [[ "${vpn_openvpn}" != true && "${vpn_wireguard}" != true ]]; then
    bashio::log.info "VPN leak monitor not started because no VPN is enabled."
    exit 0
fi

if [[ "${vpn_openvpn}" == true && "${vpn_wireguard}" == true ]]; then
    bashio::log.error "Both OpenVPN and WireGuard are enabled. Only one VPN mode is supported."
    bashio::addon.stop
    exit 1
fi

REAL_IPv4="$(read_real_ipv4)"
REAL_IPv6="$(read_real_ipv6)"

if [[ -n "${REAL_IPv4}" ]]; then
    bashio::log.info "Real (non-VPN) IPv4 from ${REAL_IP_FILEv4}: ${REAL_IPv4}"
else
    bashio::log.warning "Real IPv4 file ${REAL_IP_FILEv4} missing or empty; IP leak detection will be less strict."
fi

if [[ -n "${REAL_IPv6}" ]]; then
    bashio::log.info "Real (non-VPN) IPv6 from ${REAL_IP_FILEv6}: ${REAL_IPv6}"
else
    bashio::log.warning "Real IPv6 file ${REAL_IP_FILEv6} missing or empty; IP leak detection will be less strict."
fi

bashio::log.info "VPN detected; enabling IP leak protection and periodic monitoring."

if ! VPN_INFO_OUTv4="$(wait_for_vpn_ipv4)"; then
    bashio::log.error "Unable to obtain a VPN external IPv4 different from real IP. Stopping add-on."
    bashio::addon.stop
    exit 1
fi

if ! VPN_INFO_OUTv6="$(wait_for_vpn_ipv6)"; then
    bashio::log.error "Unable to obtain a VPN external IPv6 different from real IP. Stopping add-on."
    bashio::addon.stop
    exit 1
fi

read -r VPN_IPv4 VPN_COUNTRYv4 <<< "${VPN_INFO_OUTv4}"
read -r VPN_IPv6 VPN_COUNTRYv6 <<< "${VPN_INFO_OUTv6}"
bashio::log.info "VPN external IPv4: ${VPN_IPv4} (${VPN_COUNTRYv4}) & IPv6: ${VPN_IPv6} (${VPN_COUNTRYv6})"

if bashio::config.true 'wireguard_enabled'; then
  bashio::log.info "Requesting current status from WireGuard Client..."
#  exec wg show >/dev/stdout 2>&1 || true
  /bin/sh -c wg show > /dev/stdout
fi

while true; do
    sleep "${VPN_CHECK_INTERVAL}"

    if bashio::config.true 'wireguard_enabled'; then
      bashio::log.info "Requesting current status from WireGuard Client..."
      /bin/sh -c wg show > /dev/stdout
    fi

    if ! current_outv4="$(get_ip_infov4)"; then
        bashio::log.warning "Failed to refresh external IPv4; keeping previous assumptions."
        continue
    fi

    if ! current_outv6="$(get_ip_infov6)"; then
        bashio::log.warning "Failed to refresh external IPv6; keeping previous assumptions."
        continue
    fi

    read -r current_ipv4 current_countryv4 <<< "${current_outv4}"
    read -r current_ipv6 current_countryv6 <<< "${current_outv6}"

    bashio::log.info "Current external IP: ${current_ipv4} (${current_countryv4}) ${current_ipv6} (${current_countryv6})"

    if [[ -n "${REAL_IPv4}" ]] && [[ "${current_ipv4}" == "${REAL_IPv4}" ]]; then
        bashio::log.error "IP LEAK DETECTED: current external IPv4 ${current_ipv4} matches real IPv4 ${REAL_IPv4}. Stopping add-on."
        bashio::addon.stop
        exit 1
    fi

    if [[ -z "${REAL_IPv4}" ]] && [[ "${current_ipv4}" == "${VPN_IPv4}" ]]; then
        # No baseline to compare; only report changes in IP for visibility
        VPN_IPv4="${current_ipv4}"
    fi

    if [[ -n "${REAL_IPv6}" ]] && [[ "${current_ipv6}" == "${REAL_IPv6}" ]]; then
        bashio::log.error "IP LEAK DETECTED: current external IPv6 ${current_ipv6} matches real IPv6 ${REAL_IPv6}. Stopping add-on."
        bashio::addon.stop
        exit 1
    fi

    if [[ -z "${REAL_IPv6}" ]] && [[ "${current_ipv6}" == "${VPN_IPv6}" ]]; then
        # No baseline to compare; only report changes in IP for visibility
        VPN_IPv6="${current_ipv6}"
    fi
done
