# qbittorrent-wireguard/root/etc/s6-overlay/s6-rc.d/svc-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
#set -euo pipefail
set -e
export PATH="/usr/local/sbin:/usr/local/bin:/usr/bin:${PATH}"

# --- Configuration & Pre-checks ---

if bashio::config.true 'silent'; then
    sed -i 's|/proc/1/fd/1 hassio;|off;|g' /etc/nginx/nginx.conf
fi

# Determine log output based on silent mode
#NGINX_OUTPUT="/dev/stdout"
#if bashio::config.true 'silent'; then
#    NGINX_OUTPUT="/dev/null"
#fi

bashio::log.info "Starting nginx service"
exec nginx
#exec nginx > "${NGINX_OUTPUT}"
#exec \
#    s6-notifyoncheck -d -n 300 -w 1000 -c "nc -z localhost ${WEBUI_PORT}" \
#    s6-setuidgid abc /usr/bin/nginx > "${NGINX_OUTPUT}"
