# qbittorrent-wireguard/root/etc/s6-overlay/s6-rc.d/init-qbittorrent-config/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -e

bashio::log.info "*** init qbittorrent config ***"

##########
# INIT   #
##########

# Define preferences line
CONFIG_LOCATION=/config/qBittorrent
mkdir -p "$CONFIG_LOCATION"

# copy default config
if [ ! -f "$CONFIG_LOCATION"/qBittorrent.conf ]; then
    cp /defaults/qBittorrent.conf "$CONFIG_LOCATION"/qBittorrent.conf
fi

cd "$CONFIG_LOCATION"/ || true
LINE=$(sed -n '/\[Preferences\]/=' qBittorrent.conf) || bashio::exit.nok "qBittorrent.conf not valid"
LINE=$((LINE + 1))

# Remove unused folders
if [ -d "$CONFIG_LOCATION"/addons_config ]; then rm -r "$CONFIG_LOCATION"/addons_config; fi
if [ -d "$CONFIG_LOCATION"/qBittorrent ]; then rm -r "$CONFIG_LOCATION"/qBittorrent; fi

# Check file size
ORIGINAL_SIZE="$(wc -c "$CONFIG_LOCATION"/qBittorrent.conf)"

###########
# TIMEOUT #
###########

if bashio::config.has_value 'run_duration'; then
    echo "Timer mode set"
else
    rm -r /etc/services.d/timer
fi

##################
# Default folder #
##################

# Set configuration
if bashio::config.has_value 'SavePath'; then

    # Set variable
    DOWNLOADS=$(bashio::config 'SavePath')
    DOWNLOADS=${DOWNLOADS:-/share/downloads} # Default if not set

    # Info
    bashio::log.info "Downloads can be found in $DOWNLOADS"
fi
